# os/Darwin/Library/LaunchAgents/se.disu.mbsync.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>se.disu.mbsync.hourly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/local/bin/mbsync</string>
      <string>Amesto</string>
    </array>

    <key>Nice</key>
    <integer>1</integer>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>
