# os/Darwin/Library/LaunchAgents/se.disu.socat.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>se.disu.socat.plist</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
      <string>@SOCAT@</string>
      <string>tcp-listen:2375,reuseaddr,fork,bind=localhost</string>
      <string>unix-connect:/var/run/docker.sock</string>
    </array>
  </dict>
</plist>
