# os/Darwin/Library/LaunchAgents/se.disu.backup.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>se.disu.backup.weekly</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/now/bin/backup-home</string>
    </array>

    <key>Nice</key>
    <integer>1</integer>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Minute</key>
      <integer>0</integer>

      <key>Hour</key>
      <integer>21</integer>

      <key>Weekday</key>
      <integer>1</integer>
    </dict>
  </dict>
</plist>
